# deploy/competidor-healthcheck.service
[Unit]
Description=Verificação de saúde do CompeTIDOR
After=competidor.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/competidor.env
ExecStart=/usr/local/sbin/competidor-healthcheck
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run
